# systemd/minecraft-api.service
[Unit]
Description=Minecraft Server REST API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/minecraft/api
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python3 /home/pi/minecraft/api/server.py
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/pi/minecraft/logs/api-server.log
StandardError=append:/home/pi/minecraft/logs/api-server.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

